# roles/consul/templates/consul-client.conf.j2
; {{ ansible_managed }}

[program:consul]
command=consul agent -data-dir /tmp/consul -join={{ ansible_default_ipv4.gateway }} -config-dir=/etc/consul.d
